# friday-ddns.service
[Unit]
Description=Namecheap动态DNS更新服务
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/friday-ddns -c /etc/friday-ddns/config.yaml
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target 